# scripts/prod/supervisor/worker.conf
[program:worker] 
directory=/home/onsa/onsa/scripts/prod/ 
command=/home/onsa/onsa/scripts/prod/run_worker.sh
autostart=true 
autorestart=true 
redirect_stderr = true     
stderr_logfile=/var/log/onsa/worker.err.log
stdout_logfile=/var/log/onsa/worker.out.log 
user=onsa 